# high-vis-detection/high-vis-detection/DEBIAN/prerm
#!/bin/bash


systemctl stop high-vis-detection-app.service
systemctl disable high-vis-detection-app.service
rm -f /etc/systemd/system/high-vis-detection-app.service


systemctl daemon-reload
