# solder-quality/solder-quality/DEBIAN/prerm
#!/bin/bash


systemctl stop solder-quality-app.service
systemctl disable solder-quality-app.service
rm -f /etc/systemd/system/solder-quality-app.service


systemctl daemon-reload
