# time-in-store/time-in-store/DEBIAN/prerm
#!/bin/bash


systemctl stop time-in-store-app.service
systemctl disable time-in-store-app.service
rm -f /etc/systemd/system/time-in-store-app.service


systemctl daemon-reload
